# 80-sd-uploader.rules
SUBSYSTEM=="block", ACTION=="add", RUN+="/home/USER/.local/bin/atestscript.py"
